# .pi/systemd/embry-agent.service
[Unit]
Description=Embry Agent D-Bus Service
After=embry-session.target
PartOf=embry-session.target

[Service]
Type=simple
WorkingDirectory=%h/workspace/experiments/pi-mono
ExecStart=%h/.nvm/versions/node/v22.15.0/bin/node %h/workspace/experiments/pi-mono/packages/coding-agent/dist/dbus/cli.js --cwd %h/workspace/experiments/pi-mono
EnvironmentFile=-%h/workspace/experiments/pi-mono/.env
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PATH=%h/.nvm/versions/node/v22.15.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=EMBRY_MIN_WORKERS=2
Environment=EMBRY_MAX_WORKERS=10
Restart=always
RestartSec=5
RestartMaxDelaySec=30
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=true
ProtectHome=read-only
RestrictSUIDSGID=yes

[Install]
WantedBy=embry-session.target
